# rc.local
#!/bin/sh -e
autorunbot
/usr/sbin/mescripts
exit 0
